# systemd/openteam.target
[Unit]
Description=openteam agents target
Wants=openteam-agent@triager-01.service openteam-agent@builder-01.service openteam-agent@qa-01.service
