# systemd/openteam.target
[Unit]
Description=openteam agents target
Wants=openteam-agent@orchestrator-01.service
